# systemd/lxc.service
[Unit]
Description=Automate LXC Migration

[Service]
User=root
Type=simple
ExecStart=/bin/bash /usr/local/sbin/LXCMigration/Migration.sh
