# systemd/minute-timer.timer
[Unit]
Description=Automate LXC Migration Timer

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=multi-user.target
